# deploy/udev/99-portable-av.rules
# Trigger read-only mount helper when USB block devices appear.
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", ENV{ID_BUS}=="usb", TAG+="systemd", ENV{SYSTEMD_WANTS}="portable-av-mount@%k.service"

# Notify engine on USB removal.
ACTION=="remove", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", ENV{ID_BUS}=="usb", RUN+="/usr/local/bin/portable-av-mount --unmount --device /dev/%k"
