# deploy/systemd/portable-av-mount@.service
[Unit]
Description=Portable AV read-only mount for %i
After=local-fs.target
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/portable-av-mount --device /dev/%i
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
